# etc/cron.d/bp-addon.in
*/10 * * * * @apache_user@ @bindir@/bp-addon-session-timeout 2>&1 | logger -t bp-addon-session-timeout
